# services.d/ssh/finish
#!/usr/bin/execlineb -f
echo "Finish SSHD"
s6-svscanctl -t /var/run/s6/service
